# frc971/raspi/rootfs/frc971chrt.service
[Unit]
Description=Configure IRQs

[Service]
Type=oneshot
ExecStart=/home/pi/bin/chrt.sh

[Install]
WantedBy=multi-user.target
