# y2019/image_streamer/vision.conf
[program:vision]
command=/root/image_streamer --single_camera=true --camera0_exposure=600 --roborio_ip=10.9.71.2
redirect_stderr=false
autostart=true
autorestart=true
