# y2018/vision/exposure_loop.conf
[program:exposure_loop]
command=/root/exposure_loop.sh
redirect_stderr=false
autostart=true
autorestart=true
