# y2018/vision/vision.conf
[program:vision]
command=/root/image_streamer --single_camera=false --camera0_exposure=300 --camera1_exposure=300
redirect_stderr=false
autostart=true
autorestart=true
